# services/fcitx.service
[Unit]
Description=fcitx restarter
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/fcitx5

[Install]
WantedBy=xsession.target
